# proxy/andvari-egress-firewall.service
[Unit]
Description=Apply Andvari transparent-egress firewall rules
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service
Before=heimdall-worker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/andvari-egress-firewall apply
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
